# src/PiOTA/systemd/tunertest.service
[Unit]
Description="Test USB tuner"

[Service]
Type=simple
ExecStart=/home/tv/dtvdata/tunertest.sh




